# layers/_/etc/profile.d/99-editor.sh
# shellcheck disable=SC2148
export -- EDITOR=nvim
