# layers/nspawn/usr/local/lib/systemd/system/0-nspawn.target
[Unit]
Description = Nspawn Containers Target
