# layers/_/usr/local/lib/systemd/system/systemd-networkd.service.d/0-override.conf
[Service]
ExecStartPre = +chown -v -R -- systemd-network:systemd-network %t/systemd/network /usr/local/lib/systemd/network /usr/local/lib/systemd/networkd.conf.d
# Environment = SYSTEMD_LOG_LEVEL=debug
