# layers/_/usr/local/lib/systemd/system/systemd-networkd.service.d/0-override.conf
[Service]
ReadWritePaths = %t/systemd/network /usr/local/lib/systemd/network
ExecStartPre   = +chown -v -- systemd-network:systemd-network %t/systemd/network /usr/local/lib/systemd/network
# Environment = SYSTEMD_LOG_LEVEL=debug
