# layers/_/usr/local/lib/systemd/system/systemd-networkd.service.d/0-override.conf
[Unit]
Requires = 0-systemd-prenetworkd.service
After    = 0-systemd-prenetworkd.service
# Environment = SYSTEMD_LOG_LEVEL=debug
