# layers/_/usr/local/opt/nginx/systemd/nginx.service.d/0-override.conf
[Unit]
Upholds           = 0-nginx-sidecar.service
Requires          = 0-self-signed@%p.service
After             = 0-self-signed@%p.service

[Service]
ProtectSystem     = strict
PrivateTmp        = yes
BindReadOnlyPaths = /usr/local/opt/nginx/conf/main.nginx:%E/nginx/nginx.conf
ReadWritePaths    = %t %T %V
RuntimeDirectory  = local/%p local/%p/conf local/%p/conf/conf.d local/%p/conf/server.d local/%p/conf/http.d local/%p/conf/stream.d local/%p/conf/www
