# layers/telemetry/usr/local/lib/systemd/system/collectd.service.d/0-override.conf
[Service]
ProtectSystem     = strict
ProtectHome       = yes
PrivateTmp        = yes
ReadWritePaths    = %t/local %T %V
BindReadOnlyPaths = /usr/local/opt/%p:/%E/%p

