# layers/_/usr/local/lib/systemd/system/unattended-upgrades.service.d/0-override.conf
[Service]
KillMode       =
TimeoutStopSec = 6hours
