# layers/qemu/usr/local/lib/systemd/system/2-pemu-q35@.service.d/0-override.conf
[Unit]
PartOf      = 2-qemu-q35@%i.service

[Service]
Environment = SRC=
