# layers/router/usr/local/opt/network/1-ip-alloc@.service
[Unit]
After = 1-ip-alloc@tor.service
After = 1-ip-alloc@wg0.service
