# layers/router/usr/local/lib/systemd/system/haproxy.service.d/0-override.conf
[Service]
ProtectSystem  = strict
ProtectHome    = yes
ReadWritePaths = %t
Environment    = EXTRAOPTS="-f /usr/local/opt/haproxy/haproxy.cfg"
