# layers/router/usr/local/lib/systemd/system/haproxy.service.d/0-override.conf
[Service]
ProtectSystem     = strict
BindReadOnlyPaths = /usr/local/opt/haproxy/haproxy.cfg:%E/haproxy/haproxy.cfg
ReadWritePaths    = %t
