# layers/openvpn/usr/local/lib/systemd/system/1-vtun@.service.d/0-override.conf
[Service]
Environment = NETDEV_TEMPLATE=/usr/local/opt/openvpn/@.netdev
