# layers/router/usr/local/lib/systemd/system/ntopng.service.d/0-override.conf
[Service]
ExecStart =
ExecStart = ntopng --daemon --http-prefix=/ntop --disable-login=1
