# layers/zfs/usr/local/lib/systemd/system/zfs-mount.service.d/0-override.conf
[Service]
ExecStartPre = zfs load-key -a
