# layers/mail/usr/local/lib/systemd/system/chasquid.service.d/0-override.conf
[Unit]
Requires          =
Requires          = chasquid-smtp.socket

[Service]
Restart           = on-failure
ProtectSystem     = strict
ProtectHome       = yes
PrivateTmp        = yes

ReadWritePaths    = %t %T %V
StateDirectory    = local/%J
BindReadOnlyPaths = /usr/local/opt/%J:%E/%J
