# layers/mail/usr/local/lib/systemd/system/chasquid.service.d/0-override.conf
[Unit]
Requires          =

Requires          = 1-self-signed@%j.service
After             = 1-self-signed@%j.service
Requires          = chasquid-smtp.socket

[Service]
Restart           = on-failure
ProtectSystem     = strict
ProtectHome       = yes
PrivateTmp        = yes

ReadWritePaths    = %t %T %V
StateDirectory    = local/%J
BindReadOnlyPaths = /usr/local/opt/%J:%E/%J

ExecStartPre      = +chown -R -v chasquid:chasquid -- %C/local/self-signed/%J
