# layers/telemetry/usr/local/lib/systemd/system/systemd-journal-gatewayd.service.d/0-override.conf
[Service]
ExecStart =
ExecStart = /lib/systemd/systemd-journal-gatewayd --merge
