# layers/_/usr/local/lib/systemd/system/nftables.service.d/0-override.conf
[Unit]
RefuseManualStop = yes
Wants            = 0-tproxy-route.service
After            = 0-tproxy-route.service

[Service]
ExecReload       =
ExecStart        =
ExecStop         =

ExecStart        = /usr/local/opt/nftables/libexec/firewalld.sh
ExecReload       = /usr/local/opt/nftables/libexec/firewalld.sh
