# layers/_/usr/local/opt/mount/@.automount
[Unit]
Description    = AutoMount -- /%P

[Automount]
TimeoutIdleSec = 1h
Where          = /%P
