# layers/_/usr/local/opt/mount/@.automount
[Unit]
Description    = AutoMount -- %f

[Automount]
TimeoutIdleSec = 1h
Where          = %f
