# layers/tv/usr/local/lib/systemd/system/0-spotifyd-update.timer
[Unit]
Description        = Spotifyd Update Timer

[Timer]
OnCalendar         = weekly
RandomizedDelaySec = 1 day
