# layers/xray/usr/local/lib/systemd/system/0-xray-mon.service.d/0-override.conf
[Unit]
Description = Xray Monitor
Upholds     = 0-xray-client.service
After       = 0-xray-client.service


[Service]
Environment = PROXY_URL=socks5://127.0.0.1:50001
