# layers/telemetry/usr/local/lib/systemd/system/systemd-journal-gatewayd.socket.d/0-override.conf
[Socket]
RemoveOnStop = yes
ListenStream =
ListenStream = %t/local/systemd-journal-gatewayd.sock
