# layers/v2fly/usr/local/lib/systemd/system/dnsmasq@.service.d/2-override.conf
[Service]
Environment = UPSTREAM=127.0.0.53#52
