# layers/mail/usr/local/lib/systemd/system/opensmtpd.service.d/0-override.conf
[Service]
ProtectSystem     = strict
ProtectHome       = yes
PrivateTmp        = yes

ReadWritePaths    = %t %T %V %S/local /var/spool/smtpd
Group             = opensmtpd

BindReadOnlyPaths = /usr/local/opt/opensmtpd/smtpd.conf:%E/smtpd.conf
ExecStartPre      = /usr/local/libexec/hr-run.sh mkdir -p -- %S/local/mail/%N
ExecStartPre      = /usr/local/libexec/hr-run.sh chmod -- g+w %S/local/mail/%N
