# layers/mail/usr/local/lib/systemd/system/opensmtpd.service.d/0-override.conf
[Service]
ProtectSystem     = strict
ProtectHome       = yes
PrivateTmp        = yes

ReadWritePaths    = %t %T %V %S/local/%N /var/spool/smtpd

BindReadOnlyPaths = /usr/local/opt/opensmtpd/smtpd.conf:%E/smtpd.conf
