# layers/telemetry-collector/usr/local/lib/systemd/system/grafana-server.service.d/0-override.conf
[Service]
ProtectSystem = strict
BindPaths     = /usr/local/opt/grafana/grafana.toml:%E/grafana/grafana.ini
BindPaths     = %S/local/grafana:%S/grafana

ExecStartPost = @chmod -- o+x ${RUNTIME_DIRECTORY}
