# layers/telemetry-collector/usr/local/lib/systemd/system/grafana-server.service.d/0-override.conf
[Service]
ProtectSystem = strict
BindPaths     = %S/local/grafana:%S/grafana
