# layers/zfs/usr/local/lib/systemd/system/0-zfs-scrub-weekly@.timer
[Unit]
Description        = ZFS scrub timer -- %i

[Timer]
Persistent         = yes
OnCalendar         = %j
RandomizedDelaySec = 3 days

[Install]
WantedBy           = timers.target
