# layers/zfs/usr/local/lib/systemd/system/1-zfs-trim@.timer
[Unit]
Description        = ZFS trim timer -- %i

[Timer]
OnCalendar         = daily
RandomizedDelaySec = 12hours
