# layers/router/usr/local/lib/systemd/system/0-lan.service.d/0-override.conf
[Unit]
After       = 0-wireguard.target

[Service]
Environment = NETWORKD_DHCP=no
