# layers/telemetry-collector/usr/local/lib/systemd/system/tempo.service.d/0-override.conf
[Service]
TimeoutSec        = 2s

ProtectSystem     = strict
ProtectHome       = yes
PrivateTmp        = yes
ReadWritePaths    = %T %V

BindReadOnlyPaths = /usr/local/opt/%p/conf.yml:%E/%p/config.yml
StateDirectory    = local/%J

ExecStartPre      = +chown -v %J -- ${STATE_DIRECTORY}
