# layers/router/usr/local/lib/systemd/system/chrony.service.d/0-override.conf
[Service]
EnvironmentFile   =
BindReadOnlyPaths = /usr/local/opt/chrony/:%E/chrony

ReadWritePaths    =
ReadWritePaths    = %t %T
