# machines/esxi/usr/local/lib/systemd/system/6-cgi-youtube.sh@.service.d/0-override.conf
[Service]
CacheDirectory = adb
Environment    = "ANDROID_USER_HOME=%C/adb"
