# machines/esxi/usr/local/lib/systemd/system/6-cgi-youtube.sh@.service.d/0-override.conf
[Service]
Environment    = "HOME=%C/adb"
CacheDirectory = adb
