# layers/printer/usr/local/lib/systemd/system/cups.service.d/0-override.conf
[Unit]
ConditionPathExists = %S/local/cups/cupsd.conf

[Service]
ExecStart           =
ExecStart           = /usr/sbin/cupsd -l -c %S/local/cups/cupsd.conf
