# layers/printer/usr/local/lib/systemd/system/cups.service.d/0-override.conf
[Unit]
ConditionPathExists = /usr/local/lib/cups/cupsd.conf

[Service]
# ExecStart           =
# ExecStart           = /usr/sbin/cupsd -l -c /usr/local/lib/cups/cupsd.conf
