# layers/printer/usr/local/lib/systemd/system/cups.service.d/0-override.conf
[Service]
BindReadOnlyPaths = /usr/local/lib/cups/cupsd.conf:%E/cups/cupsd.conf
