# layers/openvpn/usr/local/lib/systemd/system/1-ovpn-tcp-server@.service.d/0-override.conf
[Unit]
Requires         = 0-ovpn-private-key.service
After            = 0-ovpn-private-key.service

BindsTo          = 1-vtun@o\x2d%i.service
After            = 1-vtun@o\x2d%i.service

[Service]
PIDFile          = %t/local/ovpn/o-%I/%J/server.pid
KillMode         =
WorkingDirectory = /
ExecStart        =

ProtectSystem    = strict

RuntimeDirectory = local/ovpn/o-%I/%J
ReadWritePaths   = %t/local/ovpn/o-%I/%J %T %V

EnvironmentFile  = %t/local/ip/%I.env
Environment      = PROTOCOL=%J
Environment      = IFACE=o-%I
Environment      = DOMAIN=o-%I.%H.home.arpa
Environment      = MARK=0x0
Environment      = CONF=%t/local/ovpn/o-%I/%J/server.ovpn

ExecStartPre     = /usr/local/libexec/envsubst2.sh /usr/local/opt/openvpn/server.ovpn ${CONF}
ExecStartPre     = dd oflag=append conv=notrunc if=/usr/local/opt/openvpn/common.ovpn of=${CONF}
ExecStart        = openvpn --config ${CONF}
