# layers/qemu/usr/local/lib/systemd/system/2-qemu-q35@.service.d/0-override.conf
[Service]
EnvironmentFile = /etc/iscsi/initiatorname.iscsi

Environment     = VNC=%t/local/qemu/%I/vnc.sock
TPM_SOCK        = %t/local/swtpm/%I/tpm.sock

ExecStartPre    =
ExecStartPre    = touch -- ${ROOT}/.live
ExecStart       =
ExecStart       = /usr/local/opt/qemu/libexec/%J@.sh --cpu ${CPU} --mem ${MEMORY} --qmp ${QMP} --monitor ${MONITOR} --bridge ${BRIDGE} --drive ${DRIVE}
