# layers/qemu/usr/local/lib/systemd/system/2-qemu-q35@.service.d/0-override.conf
[Unit]
BindTo          = 2-swtpm@%i.service
After           = 2-swtpm@%i.service

BindTo          = 2-websock-proxy@%i.socket
After           = 2-websock-proxy@%i.socket

[Service]
EnvironmentFile = /etc/iscsi/initiatorname.iscsi

Environment     = VNC_SOCK=%t/local/qemu/%I/vnc.sock
Environment     = TPM_SOCK=%t/libvirt/qemu/swtpm/%I.tpm.sock
Environment     = VIRTIO_ISO=%C/local/qemu/virtio-win.iso
Environment     = ARGV=

ExecStartPre    =
ExecStartPre    = /usr/local/opt/qemu/libexec/fs-alloc.sh %N ${FS} ${MACHINE} ${ROOT}
ExecStartPre    = /usr/local/opt/qemu/libexec/predos.sh ${MACHINE} ${ROOT} ${DRIVE}
ExecStartPre    = touch -- ${ROOT}/.live
ExecStart       =
ExecStart       = /usr/local/opt/qemu/libexec/%J@.sh \
                    --cpu ${CPU} \
                    --mem ${MEMORY} \
                    --qmp ${QMP} \
                    --monitor ${MONITOR} \
                    --tpm ${TPM_SOCK} \
                    --vnc ${VNC_SOCK} \
                    --bridge ${BRIDGE} \
                    --iscsi ${InitiatorName} \
                    --drive ${DRIVE} \
                    --drive ${VIRTIO_ISO} \
                    $ARGV
