# machines/droplet/usr/local/lib/systemd/system/var-lib-local.mount
[Unit]
Description = Mount -- %f

[Mount]
Options     = noatime
Where       = %f
What        = /dev/sda
