# machines/esxi/usr/local/lib/systemd/system/2-qemu-microvm@.service.d/0-override.conf
[Service]
Environment = MACVTAP=macvtap0
