# layers/_/usr/local/opt/nftables/nftables.service.d/0-override.conf
[Service]
ExecReload =
ExecStart  =
ExecStop   =
ExecStart  = /usr/local/opt/nftables/libexec/nft.sh
ExecReload = /usr/local/opt/nftables/libexec/nft.sh
