# layers/router/usr/local/lib/systemd/system/1-lan@.service.d/0-override.conf
[Unit]
After       = 1-ip-alloc@wg.service

[Service]
Environment = NETWORKD_DHCP=no
Environment = ENVSUBST_NETWORK=/usr/local/opt/network/lan@.network
