# layers/router/usr/local/lib/systemd/system/1-lan@.service.d/0-override.conf
[Unit]
Before       = dnsmasq@%i.service

[Service]
Environment  = NETWORKD_DHCP=no
Environment  = NETDEV=%T/_.netdev
ExecStartPre = mkdir -v -p -- ${NETWORK}.d
ExecStartPre = cp -v -f -- /usr/local/opt/network/@.network.d/0-override.conf ${NETWORK}.d/0-override.conf
ExecStartPre = chown -- systemd-network:systemd-network ${NETWORK}.d/0-override.conf
ExecStopPost = rm -v -f -- ${NETWORK}.d
