# layers/router/usr/local/lib/systemd/system/1-lan@.service.d/0-override.conf
[Unit]
After        = 1-ip-alloc@wg.service

[Service]
Environment  = NETWORKD_DHCP=no
ExecStartPre = mkdir -v -p -- ${NETWORK}.d
ExecStartPre = cp -v -f -- /usr/local/opt/network/@.network.d/0-override.conf ${NETWORK}.d/0-override.conf
ExecStopPost = rm -v -f -- ${NETWORK}.d
