# layers/router/usr/local/lib/systemd/system/squid.service.d/0-override.conf
[Service]
KillMode          =
ProtectSystem     = strict
ProtectHome       = yes
PrivateTmp        = yes
BindReadOnlyPaths = /usr/local/opt/squid/etc:%E/squid
ReadWritePaths    = %t %T %V
