# layers/telemetry-collector/usr/local/lib/systemd/system/prometheus.service.d/0-override.conf
[Service]
TimeoutStopSec  = 10s
SendSIGKILL     = yes

EnvironmentFile =
Environment     = ARGS="--config.file=/usr/local/opt/%p/%p.yml --web.external-url=/%p"

ProtectSystem   = strict
BindPaths       = %S/local/%J:%S/%J

ExecStartPre    = +chown -v %J:%J -- %S/%J
