# layers/telemetry-collector/usr/local/lib/systemd/system/prometheus.service.d/0-override.conf
[Service]
TimeoutStopSec  =
SendSIGKILL     =

EnvironmentFile =
Environment     = ARGS="--config.file=/usr/local/opt/%p/%p.yml --web.external-url=/%p"
