# machines/freenas/usr/local/lib/systemd/system/6-cgi-youtube-dl.sh@.service.d/0-override.conf
[Service]
ReadWritePaths   = %C/local /media/downloads

User             = 1000
Group            = 1000

WorkingDirectory = /media/downloads/youtube
