# layers/router/usr/local/lib/systemd/system/1-lan.service.d/0-override.conf
[Service]
Environment = NETWORKD_DHCP=no
