# machines/nuc/usr/local/lib/systemd/system/var-lib-local.mount
[Unit]
Description = Mount -- %f

[Mount]
Options     = noatime,discard
Where       = %f
What        = /dev/sdb
