# machines/pee/usr/local/opt/quadlet/3-pm-nitter@srv.container
[Unit]
Description      = Podman -- %J@%I
StopWhenUnneeded = yes

After            = network-online.target

BindsTo          = %p@redis.service
After            = %p@redis.service

[Container]
Image            = docker.io/zedeus/nitter:latest
ContainerName    = %J-%I
AutoUpdate       = registry

Label            = traefik.enable=true
Label            = traefik.http.services.%J.loadbalancer.server.port=8080
Label            = "traefik.http.routers.%J.rule=HostRegexp(`%J.{.+$$}`) || HostRegexp(`%J.{.+$$}.{.+$$}`)"

Network          = 3-pm-traefik.network

Volume           = /usr/local/opt/%J/nitter.conf:/src/nitter.conf:ro

[Service]
Restart          = on-failure

[Install]
WantedBy         = default.target
