# layers/router/usr/local/lib/systemd/system/dnsdist.service.d/0-override.conf
[Service]
ProtectSystem     = strict
BindReadOnlyPaths = /usr/local/opt/%J/%J.lua:%E/%J/%J.conf
