# layers/tv/usr/local/lib/systemd/user/gnome-remote-desktop.service.d/0-override.conf
[Unit]
Requires = 1-self-signed@%j.service
After    = 1-self-signed@%j.service
