# layers/openvpn/usr/local/lib/systemd/system/1-ip-alloc@o\x2dudp0.service.d/0-override.conf
[Unit]
After = 1-ip-alloc@o\x2dtcp0.service
