# layers/mail/usr/local/lib/systemd/system/chasquid-submission.socket.d/0-override.conf
[Socket]
ListenStream =
ListenStream = 127.0.0.53:5587
