# layers/_/usr/local/opt/cgi/06-cgi-.socket
[Unit]
Description  = CGI -- %t/cgi/%j.sock

[Socket]
RemoveOnStop = yes
Accept       = yes
ListenStream = %t/cgi/%j.sock

[Install]
WantedBy     = sockets.target
