# layers/_/usr/local/lib/systemd/system/dnsmasq@.service.d/0-override.conf
[Unit]
StopWhenUnneeded     = yes

ReloadPropagatedFrom = nftables.service
After                = nftables.service

PartOf               = sys-subsystem-net-devices-%i.device

BindsTo              = 1-ip-alloc@%i.service
After                = 1-ip-alloc@%i.service

After                = 1-vb@%i.service

[Service]
Restart              = on-failure

ProtectSystem        = strict
ProtectHome          = yes
ReadWritePaths       = %t/%J %t/local/%J

PrivateTmp           = yes
RuntimeDirectory     = local/%J/%I/conf.d local/%J/%I/hosts.d local/%J/%I/dhcp.hosts.d local/%J/%I/dhcp.opts.d

ExecStartPre         =
ExecStart            =
ExecStartPost        =
ExecStop             =

EnvironmentFile      = %t/local/ip/%I.env
EnvironmentFile      = -/usr/local/lib/etc/default/%I.%J.env
Environment          = RUN=%t/local/%J/%I
Environment          = PORT=53
Environment          = IFACE=%I
Environment          = DOMAIN=%I.%H.home.arpa
Environment          = UPSTREAM=127.0.0.53
# Environment          = NO_DHCP_IFACE=%I
Environment          = NO_DHCP_IFACE=lo
Environment          = TTL=120
Environment          = HOSTNAME=%H

ExecStartPre         = /usr/local/opt/%J/libexec/conf.sh ${RUN}/%J.conf
ExecStart            = %J --conf-file=${RUN}/%J.conf
ExecStopPost         = rm -v -fr -- ${RUN}/%J.conf

[Install]
WantedBy             =
