# layers/_/usr/local/lib/systemd/system/0-nginx-sidecar.timer
[Unit]
Description       = Reload nginx on generator updates
StopWhenUnneeded  = yes

PartOf            = nginx.service
Requisite         = nginx.service
After             = nginx.service

[Timer]
OnActiveSec       = 0
OnUnitInactiveSec = 6seconds
