# layers/_/usr/local/lib/systemd/system/0-systemd-monitor.timer
[Unit]
Description        = Polling for failed services

[Timer]
OnActiveSec        = 0
OnUnitInactiveSec  = 6minutes
RandomizedDelaySec = 2minutes
