# layers/_/usr/local/lib/systemd/system/0-wireguard.target
[Unit]
Description = Wireguard Network
Before      = default.target

[Install]
WantedBy    = default.target
