# layers/_/usr/local/lib/systemd/system/1-ns-proxy@.socket
[Unit]
Description  = netns proxy -- %t/local/ns-proxy/%J.sock

[Socket]
RemoveOnStop = yes
ListenStream = %t/local/ns-proxy/%J.sock

[Install]
WantedBy     = sockets.target
