# layers/tv/usr/local/lib/systemd/user/6-firefox@.service
[Unit]
Description = Firefox -- %I
CollectMode = inactive-or-failed

[Service]
Type        = oneshot
ExecStart   = firefox -- %I
