# layers/share/usr/local/lib/systemd/system/elasticsearch.service.d/0-override.conf
[Service]
ProtectSystem  = strict
ProtectHome    = yes
PrivateTmp     = yes

KillMode       =
SendSIGKILL    =

ReadWritePaths = %t/%J

ExecStart      =
ExecStart      = /usr/share/elasticsearch/bin/systemd-entrypoint -p ${PID_DIR}/elasticsearch.pid
