# layers/router/usr/local/lib/systemd/system/tor@default.service.d/0-override.conf
[Unit]
ReloadPropagatedFrom = nftables.service
After                = nftables.service

BindsTo              = 0-tor@%i.service
After                = 0-tor@%i.service

[Service]
ProtectSystem        = strict
BindReadOnlyPaths    = %t/local/tor/%I:%E/tor
ReadWritePaths       = %C %V %T

EnvironmentFile      = %t/local/ip/%J.env
ExecStartPost        = nft --define IPV6_NETWORK=${IPV6_IF} --file /usr/local/opt/tor/nft.conf
ExecReload           = nft --define IPV6_NETWORK=${IPV6_IF} --file /usr/local/opt/tor/nft.conf
