# layers/router/usr/local/lib/systemd/system/tor@default.service.d/0-override.conf
[Unit]
BindsTo           = 0-tor@%i.service
After             = 0-tor@%i.service

[Service]
BindReadOnlyPaths = %t/local/tor/%I:%E/tor
