# layers/share/usr/local/lib/systemd/system/6-cgi-iscsi.sh@.service.d/0-override.conf
[Service]
ProtectHome         = tmpfs
TemporaryFileSystem = /etc/rtslib-fb-target:ro
