# layers/journal/usr/local/lib/systemd/system/1-journal@.timer
[Unit]
Description        = Tail Remote Journal -- %I

[Timer]
OnActiveSec        = 0
OnUnitInactiveSec  = 2s
RandomizedDelaySec = 1s
