# layers/router/usr/local/lib/systemd/system/nftables.service.d/1-override.conf
[Unit]
Upholds    = 0-port-fwd.timer

[Service]
ExecStart  = rm -v -fr -- %t/local/port-fwd
ExecReload = rm -v -fr -- %t/local/port-fwd
