# layers/elastic-smb/usr/local/lib/systemd/system/0-smb-crawler@.timer
[Unit]
Description        = Run ElasticSearch Crawler when stopped

[Timer]
Persistent         = yes
OnCalendar         = hourly
RandomizedDelaySec = 9minutes
