# layers/rclone_mnt/usr/local/lib/systemd/system/1-rclone-mnt@.service
[Unit]
Description     = Rclone Mount -- %I
After           = network-online.target
Before          = smb.service

[Service]
Type            = notify
Restart         = on-failure
RestartSec      = 2s

# User            = 1000
# Group           = 1000
CacheDirectory  = local/rclone/%i

EnvironmentFile = /usr/local/etc/default/%i.rmount.env

Environment     = WHAT=
Environment     = ARGV=
Environment     = WHERE=/media/%I
Environment     = CACHE=%C/local/rclone/%i

ExecStartPre    = !mkdir -v -p -- /media/%I
ExecStartPre    = !chown -v -- 1000:1000 /media/%I
ExecStart       = rclone mount \
                    --allow-non-empty \
                    --allow-other \
                    --cache-dir ${CACHE} \
                    --config /dev/null \
                    --gid 1000 \
                    --uid 1000 \
                    --use-mmap \
                    --vfs-refresh \
                    --webdav-pacer-min-sleep 0 \
                    $ARGV -- ${WHAT} ${WHERE}
ExecReload      = kill -s HUP -- ${MAINPID}
ExecStop        = fusermount -u -- ${WHERE}
