# layers/certbot/usr/local/lib/systemd/system/1-certbot@.timer
[Unit]
Description        = Certbot timer -- %i

[Timer]
Persistent         = yes
OnCalendar         = daily
RandomizedDelaySec = 6hours
