# machines/freenas/usr/local/lib/systemd/system/0-bing-wallpapers.timer
[Unit]
Description        = Bing Wallpapers Timer

[Timer]
OnCalendar         = daily
RandomizedDelaySec = 12hours
