# machines/freenas/usr/local/lib/systemd/system/0-rclone-mirror.timer
[Unit]
Description        = Rclone Mirror Timer

[Timer]
OnCalendar         = daily
RandomizedDelaySec = 12hours
