# machines/freenas/usr/local/lib/systemd/system/0-rclone-update.timer
[Unit]
Description        = Rclone Update Timer

[Timer]
OnCalendar         = weekly
RandomizedDelaySec = 2days
