# machines/freenas/usr/local/lib/systemd/system/0-ytdl-update.timer
[Unit]
Description        = Youtube-DL Update Timer

[Timer]
OnCalendar         = weekly
RandomizedDelaySec = 1 day
