# machines/freenas/usr/local/lib/systemd/system/1-netns-proxy-9091@.socket
[Unit]
Description  = NetNS Proxy -- %t/local/ns-proxy/%I-%J.sock

[Socket]
RemoveOnStop = yes
ListenStream = %t/local/ns-proxy/%I-%J.sock
