# layers/_/usr/local/lib/systemd/system/nginx.service.d/0-override.conf
[Unit]
Upholds           = 0-nginx-sidecar.timer
Requires          = 1-self-signed@%p.service
After             = 1-self-signed@%p.service

[Service]
ProtectSystem     = strict
PrivateTmp        = yes
BindReadOnlyPaths = /usr/local/opt/nginx/conf/main.nginx:%E/nginx/nginx.conf
ReadWritePaths    = %t %T %V
RuntimeDirectory  = local/%p/conf/conf.d local/%p/conf/server.d local/%p/conf/http.d local/%p/conf/stream.d
