# layers/_/usr/local/lib/systemd/system/nginx.service.d/0-override.conf
[Unit]
Upholds             = 0-nginx-generator.timer
Requires            = 1-self-signed@%p.service
After               = 1-self-signed@%p.service
ConditionPathExists = /usr/local/opt/._touch

[Service]
KillMode            =
ProtectSystem       = strict
ProtectHome         = yes
PrivateTmp          = yes
ReadWritePaths      = %t %T %V

TemporaryFileSystem = %S/%J

RuntimeDirectory    = local/%p local/%p/conf/conf.d local/%p/conf/server.d local/%p/conf/http.d local/%p/conf/stream.d local/%p/conf/ssl local/%p/www
BindReadOnlyPaths   = /usr/local/opt/nginx/conf/main.nginx:%E/nginx/nginx.conf

TimeoutStopSec      = 1
