# layers/_/usr/local/lib/systemd/system/nginx.service.d/0-override.conf
[Unit]
Upholds           = 0-nginx-reload.service
Requires          = 0-self-signed@%p.service
After             = 0-self-signed@%p.service

[Service]
ProtectSystem     = strict
PrivateTmp        = yes
BindReadOnlyPaths = /usr/local/lib/nginx/nginx.conf:%E/nginx/nginx.conf
ReadWritePaths    = %t %T %V
RuntimeDirectory  = %p %p/location.d %p/http.d %p/stream.d %p/www %p/conf.d
