# layers/_/usr/local/lib/systemd/system/nginx.service.d/0-override.conf
[Unit]
Upholds           = 0-nginx-generator.timer
Requires          = 1-self-signed@%p.service
After             = 1-self-signed@%p.service

[Service]
KillMode          =
ProtectSystem     = strict
ProtectHome       = yes
PrivateTmp        = yes
ReadWritePaths    = %t %T %V
RuntimeDirectory  = local/%p/conf/conf.d local/%p/conf/server.d local/%p/conf/http.d local/%p/conf/stream.d local/%p/conf/ssl local/%p/www
BindReadOnlyPaths = /usr/local/opt/nginx/conf/main.nginx:%E/nginx/nginx.conf

ExecStopPost      = rm -fr -- %t/local/%J
TimeoutStopSec    = 1
