# layers/telemetry-collector/usr/local/lib/systemd/system/otelcol-contrib.service.d/0-override.conf
[Service]
EnvironmentFile =
# https://github.com/open-telemetry/opentelemetry-collector-contrib/tree/main/processor/transformprocessor#feature-gate
Environment     = OTELCOL_OPTIONS="--config=/usr/local/opt/otelcol-contrib/config.yml --feature-gates=transform.flatten.logs"
KillMode        =

Group           = systemd-journal
ProtectSystem   = strict
ProtectHome     = yes
PrivateTmp      = yes

StateDirectory  = local/otelcol-contrib
