# layers/cuda/usr/local/lib/systemd/system/2-nspawnd@.service.d/1-override.conf
[Service]
DeviceAllow = /dev/nvidia-modeset rwm
DeviceAllow = /dev/nvidia-uvm rwm
DeviceAllow = /dev/nvidia-uvm-tools rwm
DeviceAllow = /dev/nvidia0 rwm
DeviceAllow = /dev/nvidiactl rwm
