# layers/cuda/usr/local/lib/systemd/system/2-nspawnd@.service.d/1-override.conf
[Service]
DeviceAllow = /dev/nvidia-modeset rw
DeviceAllow = /dev/nvidia-uvm rw
DeviceAllow = /dev/nvidia-uvm-tools rw
DeviceAllow = /dev/nvidia0 rw
DeviceAllow = /dev/nvidiactl rw
