# layers/docker/usr/local/lib/systemd/system/containerd.service.d/0-override.conf
[Unit]
PartOf           = %J.slice

[Service]
Slice            = %J.slice
PropagatesStopTo = %J.slice
