# layers/docker/usr/local/lib/systemd/system/containerd.service.d/0-override.conf
[Unit]
PartOf = %j.slice

[Service]
Slice  = %j.slice
