# layers/btrfs/usr/local/lib/systemd/system/0-btrfs-scrub-weekly@.timer
[Unit]
Description        = Btrfs Scrub Timer

[Timer]
Unit               = 0-btrfs-scrub@%i.service
Persistent         = true
OnCalendar         = weekly
RandomizedDelaySec = 1d
