# layers/btrfs/usr/local/lib/systemd/system/0-btrbk.timer
[Unit]
Description        = Run btrbk every hour

[Timer]
Persistent         = yes
OnCalendar         = hourly
RandomizedDelaySec = 1m
