# layers/btrfs/usr/local/lib/systemd/system/0-btrbk.timer
[Unit]
Description        = Run btrbk every 15 minutes

[Timer]
Persistent         = yes
OnCalendar         = *:0/15
RandomizedDelaySec = 1m
