# layers/btrfs/usr/local/lib/systemd/system/1-btrbk.timer
[Unit]
Description        = Run btrbk every 15 minutes

[Timer]
OnCalendar         = 15m
RandomizedDelaySec = 1m
