# layers/tvheadend/usr/local/lib/systemd/system/tvheadend.service.d/0-override.conf
[Service]
RestartSec     = 1

ProtectSystem  = strict
ProtectHome    = yes
PrivateTmp     = yes
ReadWritePaths = %T %V
