# layers/tvheadend/usr/local/lib/systemd/system/tvheadend.service.d/0-override.conf
[Service]
RestartSec      = 1

ProtectSystem   = strict
ProtectHome     = yes
PrivateTmp      = yes
ReadWritePaths  = %t %T %V
BindPaths       = %S/local/%J:%S/%J

EnvironmentFile =
Environment     = OPTIONS="--user hts --group video --ipv6 --firstrun"
