# layers/zfs/usr/lib/systemd/system/smbd.service.d/0-override.conf
[Unit]
Requires = 0-smb-conf.service
After    = 0-smb-conf.service
