# machines/esxi/usr/local/lib/systemd/system/2-qemu-q35@.service.d/1-override.conf
[Service]
Environment = CPU_FLAGS=hypervisor=off,vme=off
