# machines/freenas/usr/local/lib/systemd/system/qbittorrent-nox@.service.d/0-override.conf
[Unit]
BindsTo              = 1-wg-netns@qbt.service
After                = 1-wg-netns@qbt.service

Wants                = 1-netns-proxy-8080@qbt.socket

RequiresMountsFor    = %S/local/qbittorrent

[Service]
Restart              = on-failure
RestartSec           = 9
ProtectSystem        = strict
ProtectHome          = yes
BindReadOnlyPaths    = %E/netns/qbt/resolv.conf:%E/resolv.conf
ReadWritePaths       = %T %V /media/downloads
StateDirectory       = local/qbittorrent

NetworkNamespacePath = %t/netns/qbt
TimeoutStopSec       = 1s

StateDirectory       = local/qbittorrent

Environment          = QBT_PROFILE=%S/local/qbittorrent
KillSignal           = SIGKILL
SuccessExitStatus    = SIGKILL
