# machines/freenas/usr/local/lib/systemd/system/qbittorrent-nox@.service.d/0-override.conf
[Unit]
# TODO:  systemd.resource-control -> RestrictNetworkInterfaces ubuntu 24
BindsTo              = 1-wg-netns@qbt.service
After                = 1-wg-netns@qbt.service

Upholds              = 1-netns-proxy-8080@qbt.socket
Before               = 1-netns-proxy-8080@qbt.socket

[Service]
Restart              = on-failure
ProtectSystem        = yes
ReadWritePaths       = %T %V %C

NetworkNamespacePath = %t/netns/qbt
TimeoutStopSec       = 30s
ExecStart            =

CacheDirectory       = local/qbt

Environment          = QBT_NO_SPLASH=TRUE
Environment          = QBT_WEBUI_LOCALHOSTAUTH=FALSE
# Environment          = QBT_PROFILE=%C
# Environment          = QBT_CONFIGURATION=%C

ExecStart            = qbittorrent-nox

# RuntimeDirectory=, StateDirectory=, CacheDirectory=, LogsDirectory=, ConfigurationDirectory=¶
# https://www.freedesktop.org/software/systemd/man/systemd.exec.html
