# machines/freenas/usr/local/lib/systemd/system/qbittorrent-nox@.service.d/0-override.conf
[Unit]
# TODO:  systemd.resource-control -> RestrictNetworkInterfaces ubuntu 24
BindsTo              = 1-wg-netns@qbt.service
After                = 1-wg-netns@qbt.service

Upholds              = 1-netns-proxy-8080@qbt.socket
Before               = 1-netns-proxy-8080@qbt.socket

RequiresMountsFor    = %S/local/qbittorrent

[Service]
Restart              = on-failure
ProtectSystem        = yes
ProtectHome          = yes
ReadWritePaths       = %T %V

NetworkNamespacePath = %t/netns/qbt
TimeoutStopSec       = 30s

StateDirectory       = local/qbittorrent

Environment          = QBT_PROFILE=%S/local/qbittorrent
