# layers/_/usr/local/lib/systemd/system/6-cgi-nginx.sh@.service.d/0-override.conf
[Service]
ProtectHome         = tmpfs
TemporaryFileSystem = %S/local/nginx
