# layers/_/usr/local/lib/systemd/system/netdata.service.d/0-override.conf
[Service]
KillMode          =
RuntimeDirectory  = local/%p
CacheDirectory    = local/%p

BindReadOnlyPaths = /usr/local/opt/%p/main.ini:%E/%p/%p.conf
BindPaths         = /usr/local/opt/%p/custom_dashboard.js:/usr/share/%p/web/custom_dashboard.js
