# layers/_/usr/local/lib/systemd/system/netdata.service.d/0-override.conf
[Service]
KillMode          =
BindReadOnlyPaths = /usr/local/opt/%p/main.conf:%E/%p/%p.conf
BindReadOnlyPaths = /usr/local/opt/%p/custom_dashboard.js:/usr/share/%p/web/custom_dashboard.js
