# layers/mail/usr/local/lib/systemd/system/apache2.service.d/0-override.conf
[Service]
KillMode            =
ProtectSystem       = strict
ProtectHome         = yes
PrivateTmp          = yes

RuntimeDirectory    = local/%J
ReadWritePaths      = %t %T %V

TemporaryFileSystem = %E/%J/sites-enabled:ro %E/%J/conf-enabled:ro
BindReadOnlyPaths   = /usr/local/opt/%J/%J.conf:%E/%J/%J.conf

BindReadOnlyPaths   = /usr/local/opt/roundcube/config.inc.php:%E/roundcube/config.inc.php

# TODO: https://forum.kopano.io/topic/3885/z-push-with-zimbra-backend-php-8-0-not-working/19
# Z-PUSH needs to be updated for php 8

BindReadOnlyPaths   = /usr/local/opt/z-push/z-push.conf.php:%E/z-push/z-push.conf.php
BindReadOnlyPaths   = /usr/local/opt/z-push/autodiscover.conf.php:%E/z-push/autodiscover.conf.php
BindReadOnlyPaths   = /usr/local/opt/z-push/imap.conf.php:%E/z-push/imap.conf.php
