# layers/mail/usr/local/lib/systemd/system/apache2.service.d/0-override.conf
[Service]
KillMode            =
ProtectSystem       = strict
ProtectHome         = yes
PrivateTmp          = yes

RuntimeDirectory    = local/%J
CacheDirectory      = local/z-push
ReadWritePaths      = %t %T %V

TemporaryFileSystem = %E/%J/sites-enabled:ro %E/%J/conf-enabled:ro
BindReadOnlyPaths   = /usr/local/opt/%J/%J.conf:%E/%J/%J.conf

BindReadOnlyPaths   = /usr/local/opt/roundcube/config.inc.php:%E/roundcube/config.inc.php

BindReadOnlyPaths   = /usr/local/opt/z-push/conf.php:/usr/share/z-push/config.php
BindReadOnlyPaths   = /usr/local/opt/z-push/autodiscover.conf.php:/usr/share/z-push/autodiscover/config.php
BindReadOnlyPaths   = /usr/local/opt/z-push/imap.conf.php:/usr/share/z-push/backend/imap/config.php

ExecStartPost       = !chown -- www-data:www-data %C/local/z-push
