# layers/postgres/usr/local/lib/systemd/system/postgresql@.service.d/0-override.conf
[Unit]
AssertPathExists    =
RequiresMountsFor   = %S/local/%p/%i
ConditionPathExists = %S/local/%p/%i

Requires            = 0-pgbouncer.service

[Service]
Restart             = always
Type                = notify-reload
PIDFile             = %t/local/%p/%i/%p.pid
RuntimeDirectory    = local/%p/%i

ProtectSystem       = strict
ProtectHome         = yes
PrivateTmp          = yes
StateDirectory      = local/%p/%i
StateDirectoryMode  = 0750
ReadWritePaths      = %t/local/%p/%i %T %V %S/local/%p/%i

User                = postgres
Group               = postgres

KillSignal          = SIGINT
TimeoutStopSec      = 1minute

Environment         = PGDATA=%S/local/%p/%i

ExecStart           =
ExecReload          =
ExecStop            =

ExecStartPre        = mkdir -v -p -- %T/pgstats
ExecStart           = /usr/local/opt/postgresql/libexec/pg-path.sh %i postgres
