# layers/postgres/usr/local/lib/systemd/system/postgresql@.service.d/0-override.conf
[Unit]
AssertPathExists    =
RequiresMountsFor   = %S/local/%p %S/local/%p/%i
ConditionPathExists = /usr/local/opt/postgresql/clusters/%i/postgresql.conf
RequiresMountsFor   = %S/local/%p/%i

[Service]
PIDFile             = %t/local/%p/%i/%p.pid
RuntimeDirectory    = local/%p/%i
StateDirectory      = local/%p/%i

ProtectSystem       = strict
ProtectHome         = yes
PrivateTmp          = yes
ReadWritePaths      = %t/local/%p/%i %T %V %S/local/%p/%i

ExecStart           =
ExecReload          =
ExecStop            =

ExecReload          = kill -s HUP -- ${MAINPID}
