# layers/postgres/usr/local/lib/systemd/system/postgresql@.service.d/0-override.conf
[Unit]
AssertPathExists    =
RequiresMountsFor   = %S/local/%p %S/local/%p/%i
ConditionPathExists = /usr/local/opt/postgresql/clusters/%i/postgresql.conf
RequiresMountsFor   = %S/local/%p/%i

[Service]
PIDFile             = %t/local/%p/%i/%p.pid
RuntimeDirectory    = local/%p/%i

ProtectSystem       = strict
ProtectHome         = yes
PrivateTmp          = yes
ReadWritePaths      = %t/local/%p/%i %T %V %S/local/%p/%i

User                = postgres
Group               = postgres

ExecStart           =
ExecReload          =
ExecStop            =

ExecStartPre        = mkdir -v -p -- %T/pgstats
ExecStart           = /usr/local/opt/postgresql/libexec/pg-path.sh %i postgres -c config_file=/usr/local/opt/postgresql/clusters/%i/postgresql.conf
ExecReload          = kill -s HUP -- ${MAINPID}
