# layers/irc/usr/local/lib/systemd/system/0-irssi.socket
[Unit]
Description  = IRSSI Socket

[Socket]
RemoveOnStop = yes
ListenStream = %t/local/irssi.sock
