# layers/router/usr/local/lib/systemd/system/nginx.service.d/1-override.conf
[Service]
BindReadOnlyPaths = %C/local/wireguard:%V/wg.sh/wireguard
