# layers/qemu/usr/local/lib/systemd/system/0-qemu.target
[Unit]
Description = Qemu VMs Target

[Install]
WantedBy    = multi-user.target
