# layers/qemu/usr/local/lib/systemd/system/0-qemu.target
[Unit]
Description       = Qemu VMs Target
RequiresMountsFor = %S/local/qemu

[Install]
WantedBy          = default.target
