# layers/qemu/usr/local/lib/systemd/system/2-swtpm@.service
[Unit]
Description      = SWTPM -- %t/local/swtpm/%I
StopWhenUnneeded = yes
CollectMode      = inactive-or-failed

[Service]
ProtectSystem    = strict
RuntimeDirectory = local/swtpm

ExecStart        = /usr/local/opt/qemu/libexec/swtpm.sh /%I %t/local/swtpm
