# layers/qemu/usr/local/lib/systemd/system/2-websock-proxy@.socket
[Unit]
Description  = Proxy qemu display socket -- %t/local/websockify/%i.sock
CollectMode  = inactive-or-failed

[Socket]
RemoveOnStop = yes
ListenStream = %t/local/websockify/%i.sock
