# layers/qemu/usr/local/lib/systemd/system/2-websock-proxy@.socket
[Unit]
Description  = Proxy qemu display socket -- %t/vm-sockets/%i.sock
CollectMode  = inactive-or-failed

[Socket]
RemoveOnStop = yes
ListenStream = %t/vm-sockets/%i.sock

[Install]
WantedBy     = sockets.target
