# layers/irc/usr/local/lib/systemd/system/znc.service.d/0-override.conf
[Service]
Restart        = on-failure

ProtectSystem  = strict
ProtectHome    = yes
ReadWritePaths = %t %T %V %S/local
User           =

Environment    = VAR=%S/local/%J
Environment    = USER=_znc

ExecStartPre   = mkdir -v -p -- ${VAR}
ExecStartPre   = chown -v -- ${USER}:${USER} ${VAR}
ExecStart      =
ExecStart      = runuser --user ${USER} -- /usr/local/opt/znc/libexec/proxy.sh znc -f --datadir ${VAR}
