# layers/irc/usr/local/lib/systemd/system/znc.service.d/0-override.conf
[Unit]
RequiresMountsFor = %S/local/%J

[Service]
Restart           = on-failure

ProtectSystem     = strict
ProtectHome       = yes
ReadWritePaths    = %t %T %V
StateDirectory    = local/%J

ExecStartPre      = !chown -v -- ${USER}:${USER} ${STATE_DIRECTORY}
ExecStart         =
ExecStart         = znc -f --datadir ${STATE_DIRECTORY}
