# machines/esxi/usr/local/lib/systemd/system/0-alarm.timer
[Unit]
Description = Alarm Timer

[Timer]
Persistent  = yes
OnCalendar  = Mon..Fri 06:30:00 America/New_York
