# layers/xray/usr/local/lib/systemd/system/0-xray-client.service.d/0-override.conf
[Unit]
BindsTo = 0-xray-server.service
After   = 0-xray-server.service

Wants   = 0-warp-mon.service
Before  = 0-warp-mon.service
PartOf  = 0-warp-mon.service

BindsTo = nginx.service
After   = nginx.service
