# layers/docker/usr/local/lib/systemd/system/0-docker.slice
[Unit]
Description = %P Slice
Before      = slices.target
