# layers/docker/usr/local/lib/systemd/system/0-docker-compose.service
[Unit]
Description     = Docker Compose

BindsTo         = docker.socket
After           = docker.socket

After           = network-online.target

[Service]
Type            = oneshot
RemainAfterExit = yes

ProtectSystem   = strict
ProtectHome     = yes
PrivateTmp      = yes
ReadWritePaths  = %t

ExecStart       = /usr/local/opt/docker-compose/libexec/up.sh
