# machines/orthanc/usr/local/lib/systemd/system/apt-daily-upgrade.service.d/1-override.conf
[Service]
Environment = NEEDRESTART_SUSPEND=1
