# layers/btrfs/usr/local/opt/mount/btrfs@.mount
[Unit]
Description = Mount -- btrfs:%f

[Mount]
Options     = noatime,compress=zstd
Where       = %f
What        = LABEL=${LABEL}
