# layers/nspawn/usr/local/lib/systemd/system/2-nspawnd@.service.d/0-override.conf
[Unit]
RequiresMountsFor = %S/local/nspawn %S/local/nspawn/%I

BindsTo           = 1-vn@nspawn.service
After             = 1-vn@nspawn.service

BindsTo           = dnsmasq@nspawn.service
After             = dnsmasq@nspawn.service

BindsTo           = 2-nspawnp@%i.service
After             = 2-nspawnp@%i.service

[Service]
DeviceAllow       = /dev/fuse rwm
ProtectSystem     = strict
ProtectHome       = yes
PrivateTmp        = yes
ReadWritePaths    = %t/systemd %T %S/local/nspawn

EnvironmentFile   = -/usr/local/etc/default/%I.nspawn.env
Environment       = ROOT=%S/local/nspawn/%I
Environment       = MACHINE=%i
Environment       = IFACE=ve-%i


ExecStartPre      = touch -- ${ROOT}/.live
ExecStopPost      = rm -v -rf -- ${ROOT}/.live

ExecStartPre      = nft -- 'add element inet user internal_ifs { ${IFACE} }'
ExecReload        = nft -- 'add element inet user internal_ifs { ${IFACE} }'
ExecStopPost      = -nft -- 'delete element inet user internal_ifs { ${IFACE} }'

ExecStart         =
ExecStart         = /usr/local/opt/nspawn/libexec/systemd-nspawn@.sh ${MACHINE} ${ROOT}
ExecStopPost      = rm -v -fr -- %t/systemd/nspawn/%i.nspawn
