# layers/nspawn/usr/local/lib/systemd/system/2-nspawnd@.service.d/0-override.conf
[Unit]
ConditionPathExists = %S/local/nspawn/%I/fs
BindsTo             = 1-vb@v\x2d%i.service
After               = 1-vb@v\x2d%i.service

BindsTo             = dnsmasq@v\x2d%i.service

Requires            = 2-nspawnp@%i.service
After               = 2-nspawnp@%i.service

PartOf              = machine.slice

[Service]
Restart             = on-failure
ProtectSystem       = strict
ProtectHome         = yes
PrivateTmp          = yes
RuntimeDirectory    = systemd/nspawn local/dnsmasq/v-%I/hosts.d
ReadWritePaths      = %t/local/dnsmasq/v-%I %t/systemd/nspawn %T %S/local/nspawn %L/journal
ReadWritePaths      = /mnt /media

EnvironmentFile     = %t/local/ip/v-%I.env
EnvironmentFile     = -/usr/local/etc/default/%I.nspawn.env

Environment         = ROOT=%S/local/nspawn/%I/fs
Environment         = MACHINE=%i
Environment         = IFACE=ve-%i
Environment         = DOMAIN=v-%I.%H.home.arpa
Environment         = HOSTS=%t/local/dnsmasq/v-%I/hosts.d

ExecStartPre        = nft --  'add    element inet user internal_ifs { ${IFACE} }'
ExecReload          = nft --  'add    element inet user internal_ifs { ${IFACE} }'
ExecStop            = -nft -- 'delete element inet user internal_ifs { ${IFACE} }'

ExecStart           =
ExecStart           = /usr/local/opt/nspawn/libexec/systemd-nspawn@.sh ${MACHINE} ${ROOT} ${HOSTS} ${DOMAIN}
