# layers/nspawn/usr/local/lib/systemd/system/2-nspawnd@.service.d/0-override.conf
[Unit]
BindsTo          = 1-vb@v\x2d%i.service
After            = 1-vb@v\x2d%i.service

Wants            = dnsmasq@v\x2d%i.service

Wants            = 2-nspawnp@%i.service
After            = 2-nspawnp@%i.service

[Service]
Restart          = on-failure
ProtectSystem    = strict
ProtectHome      = yes
PrivateTmp       = yes
RuntimeDirectory = dnsmasq/v-%I/hosts.d
ReadWritePaths   = %t/dnsmasq/v-%I %t/systemd %T %S/local/nspawn

EnvironmentFile  = %t/local/ip/v-%I.env
EnvironmentFile  = -/usr/local/etc/default/%I.nspawn.env
Environment      = ROOT=%S/local/nspawn/%I
Environment      = MACHINE=%i
Environment      = IFACE=ve-%i
Environment      = HOSTS=%t/dnsmasq/v-%I/hosts.d


ExecStartPre     = nft --  'add    element inet user internal_ifs { ${IFACE} }'
ExecReload       = nft --  'add    element inet user internal_ifs { ${IFACE} }'
ExecStop         = -nft -- 'delete element inet user internal_ifs { ${IFACE} }'

ExecStart        =
ExecStart        = /usr/local/opt/nspawn/libexec/systemd-nspawn@.sh ${MACHINE} ${ROOT} ${HOSTS}
ExecStopPost     = rm -v -fr -- %t/systemd/nspawn/%i.nspawn
