# layers/nspawn/usr/local/lib/systemd/system/2-nspawnd@.service.d/0-override.conf
[Unit]
RequiresMountsFor = %S/local/nspawn %S/local/nspawn/%I/fs

BindsTo           = 1-vn@nspawn.service
After             = 1-vn@nspawn.service

BindsTo           = 2-nspawnp@%i.service
After             = 2-nspawnp@%i.service

BindsTo           = dnsmasq@nspawn.service
After             = dnsmasq@nspawn.service

[Service]
DeviceAllow       = /dev/fuse rwm
ProtectSystem     = strict
PrivateTmp        = yes
ReadWritePaths    = %t/systemd %T %S/local/nspawn

EnvironmentFile   = -/usr/local/etc/default/%I.nspawn.env
Environment       = ROOT=%S/local/nspawn/%I/fs
Environment       = MACHINE=%i.ns
Environment       = IFACE=ve-%i.ns


ExecStartPre      = chmod +t -- ${ROOT}

ExecStartPre      = nft -- 'add element inet user internal_ifs { ${IFACE} }'
ExecReload        = nft -- 'add element inet user internal_ifs { ${IFACE} }'
ExecStopPost      = -nft -- 'delete element inet user internal_ifs { ${IFACE} }'


ExecStart         =
ExecStart         = /usr/local/opt/nspawn/libexec/systemd-nspawn@.sh ${MACHINE} ${ROOT}
ExecStartPost     = /usr/local/libexec/retry.sh 0.1 ping -c 1 -w 1 -- ${MACHINE}

ExecStopPost      = chmod -t -- ${ROOT}
