# layers/v2fly/usr/local/lib/systemd/system/0-v2update.timer
[Unit]
Description        = V2fly Update Timer

[Timer]
OnCalendar         = weekly
RandomizedDelaySec = 1 day
