# layers/podman/usr/local/lib/systemd/system/0-podman.target
[Unit]
Description = Podman Containers Target

After       = network-online.target

[Install]
WantedBy    = multi-user.target
