# layers/podman/usr/local/lib/systemd/system/0-podman.target
[Unit]
Description       = Podman Containers Target
RequiresMountsFor = %S/local/docker

After             = network-online.target

[Install]
WantedBy          = default.target
