# layers/podman/usr/local/lib/systemd/system/0-restart-failed-pods.timer
[Unit]
Description       = Restart failed podman services

[Timer]
OnActiveSec       = 12minutes
OnUnitInactiveSec = 12minutes
