# layers/podman/usr/local/lib/systemd/system/3-loopback-6to4@.socket
[Unit]
Description      = Socket 6to4 -- %I:localhost:%I
StopWhenUnneeded = yes

[Socket]
BindIPv6Only     = ipv6-only
ListenStream     = %I
