# layers/gh_mirror/usr/local/lib/systemd/system/1-gh-mirror@.timer
[Unit]
Description        = Github Mirror Timer -- %I

[Timer]
Persistent         = yes
OnCalendar         = daily
RandomizedDelaySec = 12hours
