# layers/oci-registry/usr/local/lib/systemd/system/1-registry-gc@.service.d/0-override.conf
[Unit]
Description      = OCI Registry GC Service -- %I

BindsTo          = 1-registry-proxy@%i.service
After            = 1-registry-proxy@%i.service

[Service]
RuntimeDirectory =
RuntimeDirectory = local/oci-registry/gc/%I
Environment      = CONF=%t/local/oci-registry/gc/%I/conf.yml

ExecStart        =
ExecStart        = registry garbage-collect ${CONF}
