# layers/zfs/usr/local/lib/systemd/system/zfs-scrub@.service.d/0-override.conf
[Unit]
ConditionPathExists = /usr/local/opt/zfs/libexec/scrub.sh

[Service]
Type                = oneshot
ExecStart           =
ExecStart           = /usr/local/opt/zfs/libexec/scrub.sh %I
ExecStop            =
ExecStop            = zpool scrub -p -- %I
