# layers/share/usr/local/lib/systemd/system/0-wsdd.service
[Unit]
Description   = WSDD

[Service]
Restart       = on-failure
ProtectSystem = strict
ProtectHome   = yes
PrivateTmp    = yes

ExecStart     = wsdd
