# layers/share/usr/local/lib/systemd/system/1-fs-veto@.timer
[Unit]
Description        = FS Veto Timer -- %f

[Timer]
OnCalendar         = weekly
RandomizedDelaySec = 1 day
