# machines/freenas/usr/local/lib/systemd/system/gerbera.service.d/0-override.conf
[Service]
ProtectSystem       = yes
ProtectHome         = yes
TemporaryFileSystem = /media:ro /mnt:ro /var:ro /usr/local:ro

BindReadOnlyPaths   = /usr/local/opt/%N:%E/%N
