# machines/freenas/usr/local/lib/systemd/system/gerbera.service.d/0-override.conf
[Service]
ProtectSystem       = yes
ProtectHome         = yes
BindReadOnlyPaths   = /usr/local/opt/%N:%E/%N
TemporaryFileSystem = /mnt:ro /media:ro
