# layers/zfs/usr/local/lib/systemd/system/smbd.service.d/1-override.conf
[Unit]
Wants        = 0-smb-share.service

[Service]
ExecStartPre = systemd-tmpfiles --create --remove -- %E/tmpfiles.d/0-smbd.conf
