# layers/telemetry-collector/usr/local/lib/systemd/system/loki.service.d/0-override.conf
[Service]
TimeoutSec     =

ProtectSystem  = strict
ProtectHome    = yes
PrivateTmp     = yes
ReadWritePaths = %T %V

BindPaths      = /usr/local/opt/%p/conf.yml:%E/%p/config.yml
