# layers/telemetry-collector/usr/local/lib/systemd/system/loki.service.d/0-override.conf
[Service]
TimeoutSec     = 2s

ProtectSystem  = strict
ProtectHome    = yes
PrivateTmp     = yes
ReadWritePaths = %T %V %S/%J

BindPaths      = /usr/local/opt/%p/conf.yml:%E/%p/config.yml
BindPaths      = %S/local/%J:%S/%J

