# layers/router/usr/local/lib/systemd/system/0-dhcpd.timer
[Unit]
Description       = TelusIA_PD timer

[Timer]
OnStartupSec      = 6seconds
OnUnitInactiveSec = 60seconds
