# layers/wireguard/usr/local/lib/systemd/system/1-wg@.service.d/0-override.conf
[Unit]
Requires         = 0-wg-private-key.service
After            = 0-wg-private-key.service

[Service]
RuntimeDirectory = local/wireguard/dnsmasq.d
StateDirectory   = local/wireguard
CacheDirectory   = local/wireguard

Environment      = NETDEV=%T/_.netdev
Environment      = WG_NETDEV=%t/systemd/network/0-%I.netdev

Environment      = HOSTNAME=%H
Environment      = WG_SERVER_NAME=
Environment      = WG_PEERS=
Environment      = DNSMASQD=%t/local/wireguard/dnsmasq.d/%I.hosts
Environment      = STATE=%S/local/wireguard
Environment      = CACHE=%C/local/wireguard

ExecStartPre     = /usr/local/opt/wireguard/libexec/alloc.sh ${IFACE} ${DOMAIN} ${WG_NETDEV} ${DNSMASQD} ${CACHE} ${STATE} ${WG_SERVER_NAME} ${WG_PEERS}
ExecStartPre     = chown -v -- systemd-network:systemd-network ${WG_NETDEV}
