# layers/_/usr/local/lib/systemd/system/haproxy.service.d/0-override.conf
[Unit]
Requires          = 1-self-signed@%p.service
After             = 1-self-signed@%p.service

[Service]
ProtectSystem     = strict
ProtectHome       = yes
ReadWritePaths    = %t

RuntimeDirectory  = %p

Environment       = EXTRAOPTS="-f /usr/local/opt/haproxy/gen.cfg"
BindReadOnlyPaths = /usr/local/opt/%J/%J.cfg:%E/%J/%J.cfg

ExecStartPre      = ln -snf -- %C/local/self-signed/%p/ssl.pem %t/%p/ssl.crt
