# layers/zfs/usr/lib/systemd/system/0-zfs-trim-daily@.timer
[Unit]
Description        = ZFS trim timer -- %i

[Timer]
OnStartupSec       = %j
OnUnitInactiveSec  = %j
RandomizedDelaySec = 12 hours

[Install]
WantedBy           = timers.target
