# layers/router/usr/local/lib/systemd/system/dnsmasq@.service.d/0-override.conf
[Unit]
StopWhenUnneeded     = yes
CollectMode          = inactive-or-failed

ReloadPropagatedFrom = nftables.service
After                = nftables.service

PartOf               = sys-subsystem-net-devices-%i.device

BindsTo              = 0-ip-alloc@%i.service
After                = 0-ip-alloc@%i.service
BindsTo              = 0-ip-alloc@vn\x2d%i.service
After                = 0-ip-alloc@vn\x2d%i.service

[Service]
Restart              = on-failure

ProtectSystem        = strict
ReadWritePaths       = %t/dnsmasq %t/local/dnsmasq

PrivateTmp           = yes
RuntimeDirectory     = local/dnsmasq/%I local/dnsmasq/%I/conf.d local/dnsmasq/%I/hosts.d local/dnsmasq/%I/dhcp.hosts.d local/dnsmasq/%I/dhcp.opts.d

ExecStartPre         =
ExecStart            =
ExecStartPost        =
ExecStop             =

EnvironmentFile      = %t/local/ip/%I.env
EnvironmentFile      = -/usr/local/lib/etc/default/%I.dnsmasq.env
Environment          = RUN=%t/local/dnsmasq/%I
Environment          = IFACE=%I
Environment          = DOMAIN=%I
Environment          = TTL=120
Environment          = HOST=%H

ExecStartPre         = /usr/local/libexec/envsubst2.sh /usr/local/opt/dnsmasq/dnsmasq.conf ${RUN}/dnsmasq.conf
ExecStart            = dnsmasq --conf-file=${RUN}/dnsmasq.conf

[Install]
WantedBy             =
