# layers/router/usr/local/lib/systemd/system/dnsmasq@.service.d/0-override.conf
[Unit]
StopWhenUnneeded     = yes
CollectMode          = inactive-or-failed

ReloadPropagatedFrom = nftables.service
After                = nftables.service

PartOf               = sys-subsystem-net-devices-%i.device

BindsTo              = 0-ip-alloc@%i.service
After                = 0-ip-alloc@%i.service
BindsTo              = 0-ip-alloc@vn\x2d%i.service
After                = 0-ip-alloc@vn\x2d%i.service

[Service]
Restart              = on-failure

ProtectSystem        = strict
ReadWritePaths       = %t/dnsmasq

PrivateTmp           = yes
RuntimeDirectory     = dnsmasq/%I dnsmasq/%I/conf.d dnsmasq/%I/hosts.d dnsmasq/%I/dhcp.hosts.d dnsmasq/%I/dhcp.opts.d
EnvironmentFile      = -/usr/local/lib/env.d/dnsmasq.%I.env

ExecStartPre         =
ExecStart            =
ExecStartPost        =
ExecStop             =

Environment          = INSTANCE=%I
Environment          = IFACE=%I
Environment          = DOMAIN=%I
Environment          = TTL=120
Environment          = HOST=%H

ExecStartPre         = envsubst2.sh /usr/local/lib/dnsmasq/5-if.conf %t/dnsmasq/%I/conf.d/5-if.conf
ExecStartPre         = /usr/local/lib/dnsmasq/ip.sh ${INSTANCE} ${IFACE} ${DOMAIN} ${TTL}
ExecStartPre         = nft -- 'add element inet user internal_ifs { ${IFACE} }'
ExecReload           = nft -- 'add element inet user internal_ifs { ${IFACE} }'
ExecStart            = dnsmasq --conf-dir=/usr/local/lib/dnsmasq/conf.d,*.conf --conf-dir=%t/dnsmasq/%I/conf.d,*.conf
ExecStopPost         = -nft -- 'delete element inet user internal_ifs { ${IFACE} }'

[Install]
WantedBy             =
