# layers/router/usr/local/lib/systemd/system/dnsmasq@.service.d/1-override.conf
[Service]
Environment = NO_DHCP_IFACE=lo
Environment = TTL=600
