# machines/droplet/usr/local/lib/systemd/system/0-k3s.service.d/0-override.conf
[Unit]
After = haproxy.service
After = nginx.service
