# layers/gateway/usr/local/lib/systemd/system/nginx.service.d/1-override.conf
[Unit]
Wants            = 0-htpasswd.service

[Service]
RuntimeDirectory = local/nginx/acme

ExecStartPre     = chown www-data:www-data -- %t/local/nginx/acme
