# layers/gateway/usr/local/lib/systemd/system/nginx.service.d/1-override.conf
[Unit]
Wants             = 0-htpasswd.service
RequiresMountsFor = %S/local/htpasswd

[Service]
RuntimeDirectory  = local/nginx/acme

ExecStartPre      = mkdir -p -v -- %t/local/nginx/acme/.well-known
ExecStartPre      = chown --recursive www-data:www-data -- %t/local/nginx/acme
