# layers/gateway/usr/local/lib/systemd/system/nginx.service.d/1-override.conf
[Unit]
RequiresMountsFor = %S/local/htpasswd

[Service]
RuntimeDirectory  = local/nginx/acme

Environment       = HTPASSWD_SECRET=%H*%m
EnvironmentFile   = -/usr/local/etc/default/htpasswd.env

ExecStartPre      = mkdir -p -v -- %t/local/nginx/acme/.well-known
ExecStartPre      = chown --recursive www-data:www-data -- %t/local/nginx/acme
