# layers/nspawn/usr/local/lib/systemd/system/1-nspawnd@.service.d/0-override.conf
[Unit]
BindsTo         = 0-vn@nspawn.service
After           = 0-vn@nspawn.service

BindsTo         = 1-nspawnp@%i.service
After           = 1-nspawnp@%i.service

[Service]
DeviceAllow     = /dev/fuse rwm
ProtectSystem   = strict
PrivateTmp      = yes
ReadWritePaths  = %t/systemd %T %S/local/nspawn

EnvironmentFile = -/usr/local/etc/default/%I.nspawn.env
Environment     = ROOT=%S/local/nspawn/%I/fs
Environment     = MACHINE=%i.nspawn

ExecStartPre    = chmod +t -- ${ROOT}
ExecStart       =
ExecStart       = systemd-nspawn --keep-unit --settings=override --boot --notify-read=yes -U --network-veth --link-journal=try-guest --machine=${MACHINE} --directory=${ROOT}
ExecStartPost   = /usr/local/opt/nspawn/libexec/start-post.sh ${MACHINE}
ExecStopPost    = chmod -t -- ${ROOT}
