# layers/nspawn/usr/local/lib/systemd/system/1-nspawnd@.service.d/0-override.conf
[Unit]
BindsTo         = 0-vn@nspawn.service
After           = 0-vn@nspawn.service

[Service]
DeviceAllow     = /dev/fuse rwm
ProtectSystem   = strict
PrivateTmp      = yes
ReadWritePaths  = %t/systemd %S/local/machines

EnvironmentFile = -/usr/local/etc/default/%I.nspawn.env

Environment     = ROOT=%S/local/machines/%I

ExecStart       =
ExecStart       = systemd-nspawn --keep-unit -U --settings=override --machine=%i --directory=${ROOT}/fs

ExecStartPre    = /usr/local/opt/nspawn/libexec/start-pre.sh %I ${ROOT}/fs
ExecStartPre    = chmod +t -- ${ROOT}
ExecStartPost   = /usr/local/opt/nspawn/libexec/start-post.sh %I
ExecStopPost    = chmod -t -- ${ROOT}
