# layers/nspawn/usr/local/lib/systemd/system/1-nspawnd@.service.d/0-override.conf
[Unit]
BindsTo         = 0-vn@nspawn.service
After           = 0-vn@nspawn.service
BindsTo         = dnsmasq@nspawn.service

[Service]
DeviceAllow     = /dev/fuse rwm
ProtectSystem   = strict
PrivateTmp      = yes
ReadWritePaths  = %t/systemd %S/local/machines

EnvironmentFile = -/usr/local/etc/default/%I.nspawn.env

Environment     = FS=%S/local/machines/%I
Environment     = TMPFILES_EXTRA=%T/tmpfiles.extra

ExecStart       =

ExecStartPre    = touch -- ${TMP_EXTRA}
ExecStartPre    = chmod +t -- ${FS}

# TODO: check if ubuntu:24 can use  LoadCredential = ...
ExecStart       = systemd-nspawn --quiet --keep-unit --boot --link-journal=try-guest --network-veth -U --settings=override \
                   --load-credential ssh.authorized_keys.root:%h/.ssh/authorized_keys \
                   --load-credential tmpfiles.extra:${TMP_EXTRA}

ExecStartPost   = /usr/local/opt/nspawn/libexec/wait.sh %I
ExecStopPost    = chmod -t -- ${FS}
