# layers/nspawn/usr/local/lib/systemd/system/1-nspawnd@.service.d/0-override.conf
[Unit]
BindsTo         = 0-vn@nspawn.service
After           = 0-vn@nspawn.service

BindsTo         = 1-nspawnp@%i.service
After           = 1-nspawnp@%i.service

[Service]
DeviceAllow     = /dev/fuse rwm
ProtectSystem   = strict
PrivateTmp      = yes
ReadWritePaths  = %t/systemd %T %S/local/nspawn

EnvironmentFile = -/usr/local/etc/default/%I.nspawn.env
Environment     = ROOT=%S/local/nspawn/%I/fs

ExecStartPre    = chmod +t -- ${ROOT}
ExecStart       =
ExecStart       = systemd-nspawn --keep-unit -U --settings=override --machine=%i --directory=${ROOT}
ExecStopPost    = chmod -t -- ${ROOT}
