# layers/qemu/usr/local/lib/systemd/system/2-pemu-microvm@.service.d/0-override.conf
[Unit]
AssertPathExists = %C/local/qemu/cloud.img/raw
