# layers/_/usr/local/lib/systemd/system/otelcol-contrib.service.d/0-override.conf
[Service]
EnvironmentFile =
Environment     = OTELCOL_OPTIONS="--config=/usr/local/opt/otelcol-contrib/config.yml"
KillMode        =

ProtectSystem   = strict
ProtectHome     = yes
PrivateTmp      = yes
