# layers/_/usr/local/lib/systemd/system/apt-daily-upgrade.service.d/0-override.conf
[Service]
KillMode       =
TimeoutStopSec = 10minutes
