# layers/media/usr/local/lib/systemd/system/plexmediaserver.service.d/0-override.conf
[Unit]
RequiresMountsFor   = %S/local/%J

[Service]
ProtectSystem       = strict
ProtectHome         = yes
PrivateTmp          = yes
ReadWritePaths      = %T %V
StateDirectory      = local/%J

User                = 1000
Group               = 1000

TemporaryFileSystem = /usr/local:ro %E:ro /var:ro /mnt:ro /srv:ro
Environment         = LC_ALL=zh_CN.UTF-8
