# layers/k3s/usr/local/lib/systemd/system/1-ip-alloc@k3s.service.d/0-override.conf
[Unit]
After = 1-ip-alloc@o\x2dtcp0.service
After = 1-ip-alloc@o\x2dudp0.service
After = 1-ip-alloc@w\x2d0.service
