# layers/_/usr/local/lib/systemd/system/ttyd.service.d/0-override.conf
[Service]
Environment   = SOCK=%t/ttyd.sock
ExecStart     =
ExecStart     = ttyd --interface ${SOCK} --base-path /ttyd -- bash
# TODO: use actual isolated nginx user, use --credential flag later
ExecStartPost = /usr/local/libexec/retry.sh 0.1 stat -- ${SOCK}
ExecStartPost = chown -- nobody ${SOCK}
ExecStopPost  = rm -v --recursive --force -- ${SOCK}
