# layers/_/usr/local/lib/systemd/system/ttyd.service.d/0-override.conf
[Service]
EnvironmentFile = /usr/local/opt/ttyd/options.env
Environment     = TTYD_SOCK=%t/local/ttyd.sock
Environment     = HOME=%h
Environment     = USER=%u
Environment     = SHELL=/usr/bin/bash
Environment     = HOSTNAME=%H

ExecStart       =
ExecStart       = env --unset INVOCATION_ID -- ttyd --interface ${TTYD_SOCK} --socket-owner www-data:www-data --base-path /ttyd --ping-interval 6 -t disableLeaveAlert=true -t fontSize=14 -t theme=${TTYD_THEME} --writable --cwd ${HOME} -- ${SHELL} --login
ExecStopPost    = rm -v --recursive --force -- ${TTYD_SOCK}
