# layers/_/usr/local/lib/systemd/system/ttyd.service.d/0-override.conf
[Service]
EnvironmentFile = /usr/local/opt/ttyd/options.env
Environment     = TTYD_SOCK=%t/ttyd.sock
Environment     = HOME=%h
Environment     = USER=%u
Environment     = HOSTNAME=%H

ExecStart       =
ExecStart       = ttyd --interface ${TTYD_SOCK} --base-path /ttyd --client-option theme=${TTYD_THEME} --cwd ${HOME} -- bash
# TODO: use actual isolated nginx user, use --socket-owner flag later
ExecStartPost   = /usr/local/libexec/retry.sh 0.1 stat -- ${TTYD_SOCK}
ExecStartPost   = chown -- nobody ${TTYD_SOCK}
ExecStopPost    = rm -v --recursive --force -- ${TTYD_SOCK}
