# machines/esxi/usr/local/lib/systemd/system/2-nspawnd@.service.d/1-override.conf
[Service]
DeviceAllow = /dev/fuse rwm
DeviceAllow = /dev/nvidia0 rwm
DeviceAllow = /dev/nvidiactl rwm
