# layers/_/usr/local/lib/systemd/system/systemd-resolved.service.d/0-override.conf
[Service]
ExecStartPre = +-chown -v -R -- systemd-resolve:systemd-resolve /usr/local/lib/systemd/dnssd /usr/local/lib/systemd/resolved.conf.d
