# layers/gateway/usr/local/lib/systemd/system/0-certbot-update.timer
[Unit]
Description        = Certbot update timer

[Timer]
Persistent         = yes
OnCalendar         = weekly
RandomizedDelaySec = 2days
