# layers/ddns/usr/local/lib/systemd/system/0-ddclient@.timer
[Unit]
Description        = DDNS Timer - %I

[Timer]
OnStartupSec       = 6seconds
OnUnitInactiveSec  = 10minutes
RandomizedDelaySec = 1minute
