# layers/openvpn/usr/local/lib/systemd/system/1-ovpn-tcp@.service.d/0-override.conf
[Unit]
Requires         = 1-self-signed@ovpn.service
After            = 1-self-signed@ovpn.service

Requires         = 0-ovpn-gen-keys.service
After            = 0-ovpn-gen-keys.service

BindsTo          = 1-vtun@o\x2d%j%i.service
After            = 1-vtun@o\x2d%j%i.service

[Service]
PIDFile          = %t/local/ovpn/%J%I/pid
KillMode         =
WorkingDirectory = /
ExecStart        =

ProtectSystem    = strict

RuntimeDirectory = local/ovpn/%J%I
ReadWritePaths   = %t/local/ovpn/%J%I %T %V

EnvironmentFile  = %t/local/ip/o-%J%I.env
Environment      = HOSTNAME=%H
Environment      = RUN=%t/local/ovpn/%J%I
Environment      = SSL=%C/local/self-signed/ovpn
Environment      = TLS_CRYPT=%S/local/openvpn/tls-crypt-v2-server.key
Environment      = PROTOCOL=%J
Environment      = IFACE=o-%J%I
Environment      = MARK=0x0

ExecStartPre     = /usr/local/libexec/envsubst2.sh /usr/local/opt/openvpn/server.ovpn ${RUN}/server.ovpn
ExecStartPre     = dd oflag=append conv=notrunc if=/usr/local/opt/openvpn/common.ovpn of=${RUN}/server.ovpn
ExecStartPre     = dd oflag=append conv=notrunc if=/usr/local/opt/openvpn/%J.server.ovpn of=${RUN}/server.ovpn
ExecStart        = openvpn --config ${RUN}/server.ovpn
