# layers/openvpn/usr/local/lib/systemd/system/1-ovpn-tcp@.service.d/0-override.conf
[Unit]
Requires         = 0-ovpn-private-key.service
After            = 0-ovpn-private-key.service

Requires         = 1-self-signed@ovpn\x2do\x2d%j%i.service
After            = 1-self-signed@ovpn\x2do\x2d%j%i.service

BindsTo          = 1-vtun@o\x2d%j%i.service
After            = 1-vtun@o\x2d%j%i.service

[Service]
PIDFile          = %t/local/ovpn/o-%I/%J/server.pid
KillMode         =
WorkingDirectory = /
ExecStart        =

ProtectSystem    = strict

RuntimeDirectory = local/ovpn/o-%I/%J
ReadWritePaths   = %t/local/ovpn/o-%I/%J %T %V

EnvironmentFile  = %t/local/ip/o-%I.env
Environment      = RUN=%t/local/ovpn/o-%I/%J
Environment      = PROTOCOL=%J
Environment      = IFACE=o-%I
Environment      = HOME_DOMAIN=%H.home.arpa
Environment      = MARK=0x0
Environment      = CONF=%t/local/ovpn/o-%I/%J/server.ovpn

ExecStartPre     = /usr/local/libexec/envsubst2.sh /usr/local/opt/openvpn/server.ovpn ${RUN}/server.ovpn
ExecStartPre     = dd oflag=append conv=notrunc if=/usr/local/opt/openvpn/common.ovpn of=${RUN}/server.ovpn
ExecStart        = openvpn --config ${RUN}/server.ovpn
