# layers/openvpn/usr/local/lib/systemd/system/1-ovpn-tcp@.service.d/1-override.conf
[Service]
Environment = PROTOCOL=%J6-server
