# layers/cloudflare-warp/usr/local/lib/systemd/system/0-privoxy@.service.d/0-override.conf
[Unit]
ConditionPathExists = /usr/local/opt/%J/conf.d/%i.conf

[Service]
ProtectSystem       = strict
ProtectHome         = yes
PrivateTmp          = yes

ReadWritePaths      = %t %T %V
BindReadOnlyPaths   = /usr/local/opt/%J/conf.d/%i.conf:%E/%J/%i.conf

Environment         = PIDFILE=%t/%J-%i.pid
Environment         = CONFIGFILE=%E/%J/%i.conf
