# layers/cloudflare-warp/usr/local/lib/systemd/system/warp-svc.service.d/0-override.conf
[Service]
ProtectSystem  = strict
ProtectHome    = yes
PrivateTmp     = yes

StateDirectory =
LogsDirectory  =

BindPaths      = %S/local/warp:%S/cloudflare-warp
ExecStartPost  = /usr/local/libexec/retry.sh 0.1 stat -- %t/cloudflare-warp/warp_service
ExecStartPost  = warp-cli --accept-tos mode proxy
ExecStartPost  = warp-cli --accept-tos connect
