# layers/_/usr/local/lib/systemd/system/avahi-daemon.service.d/0-override.conf
[Service]
ExecStart =
ExecStart = %N --syslog --file /usr/local/opt/avahi/avahi-daemon.conf
